# conf/supervisor/gstfsd.conf
[program:gstfsd]
command=/usr/bin/python /usr/local/bin/gstfsd
directory=/usr/local/bin
user=root
autostart=true
autorestart=true
redirect_stderr=true